# src/supervisord.conf
[supervisord]
nodaemon=false

[program:cassandra]
user=cassandra
command=cassandra -f